# site.conf
[program:site]
user=root
command=/uwsgi/uwsgi --ini /uwsgi/uwsgi.ini
directory=/site
stopsignal=QUIT
stdout_logfile=/tmp/site.stdout.log
stderr_logfile=/tmp/site.stderr.log
